# etc/profile.d/mouse_init.sh
#!/usr/bin/bash
synclient VertScrollDelta=-69
synclient HorizScrollDelta=-69

